# systemd/alarmmonitor-network-portal.service
[Unit]
Description=Alarmmonitor WLAN Konfigurationsportal
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/alarmmonitor-network-portal
ExecStart=/opt/alarmmonitor-network-portal/venv/bin/python network_portal.py
Restart=on-failure
Environment=WIFI_INTERFACE=wlan0
Environment=PORTAL_PORT=8080
Environment=WPA_SUPPLICANT_CONF=/etc/wpa_supplicant/wpa_supplicant.conf
Environment=HOTSPOT_SUBNET=10.20.0.0/24
Environment=PORTAL_SECRET=change-me
User=root
Group=root

[Install]
WantedBy=multi-user.target
